# cron.d/squert-ip2c
# /etc/cron.d/squert-ip2c
#
# crontab entry to update SQueRT ip2c mappings

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/5 * * * * root /usr/bin/so-squert-ip2c-5min > /dev/null 2>&1
00 6 * * 0 root /usr/bin/so-squert-ip2c cron > /dev/null 2>&1
